# systemd/game-night.socket
[Unit]
Description=Game Night Socket
After=network-online.target

[Socket]
ListenStream=8000

[Install]
WantedBy=sockets.target
